# zoa-tester.service
[Unit]
Description=ZOA Rate Limit Tester (uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/api-test
ExecStart=/usr/local/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
